# digitalocean/packer/files/etc/systemd/system/shinyproxy.service.d/override.conf
[Service]
EnvironmentFile=/etc/shinyproxy/.env
